# deploy/despatch-updater.path
[Unit]
Description=Watch for Despatch update requests

[Path]
PathExists=/var/lib/despatch/update/request/update-request.json
PathExistsGlob=/var/lib/despatch/update/request/update-request-*.json
Unit=despatch-updater.service

[Install]
WantedBy=multi-user.target
