# deploy/despatch-pam-reset-helper.service
[Unit]
Description=Despatch PAM Reset Helper
After=despatch-pam-reset-helper.socket
Requires=despatch-pam-reset-helper.socket

[Service]
Type=simple
WorkingDirectory=/opt/despatch
EnvironmentFile=/opt/despatch/.env
ExecStart=/opt/despatch/despatch-pam-reset-helper
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
PrivateNetwork=true
RestrictAddressFamilies=AF_UNIX
ProtectHome=true
# Keep /etc writable for chpasswd integration.
ProtectSystem=off
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
